# posten-levering/app/crontab
30 4 * * * python /app/main.py $POSTALCODE $MQTT_HOST $MQTT_PORT $MQTT_USER $MQTT_PASS